# weather.timer
[Unit]
Description=Run weather report every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
